# config/autostart/dropbox.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=dropbox.py start
Name=Dropbox
